# ngrok-watcher.service
[Unit]
Description=ngrok watcher

[Service]
ExecStart=/home/aghinassi/ngrok-watcher/main.sh

[Install]
WantedBy=default.target
